# images/nova-qemu/supervisord.conf
[unix_http_server]
file = /var/run/supervisor.sock

[supervisord]
nodaemon=true
logfile = /var/log/supervisor/supervisord.log
logfile_maxbytes = 200KB
logfile_backups = 1
pidfile = /var/run/supervisord.pid
childlogdir = /var/log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

[program:libvirtd]
command:/usr/sbin/libvirtd
priority=10
startsec=10
stderr_events_enabled=true
stdout_events_enabled=true

[program:virtlogd]
command:/usr/sbin/virtlogd
priority=10
startsec=10
stderr_events_enabled=true
stdout_events_enabled=true

[program:nova-api-metadata]
command:/usr/sbin/nova-api-metadata --log-file=/var/log/nova/nova-metadata.log
priority=10
startsec=10
stderr_events_enabled=true
stdout_events_enabled=true

[program:nova-compute]
command:/usr/sbin/ova-compute --config-file /etc/nova/nova.conf --config-file /etc/nova/nova-compute.conf --log-file=/var/log/nova/nova-compute.log
priority=10
startsec=10
stderr_events_enabled=true
stdout_events_enabled=true
